# client/crontab
3 * * * * /var/AutoProxy/client/request.sh /var/AutoProxy/client/request.conf >> /var/AutoProxy/client/request.log
